# scripts/fd-rdd.service
[Unit]
Description=fd-rdd: RDD-based elastic file indexer (user service)
After=local-fs.target
Documentation=https://github.com/AegisCompute/fd-rdd

[Service]
Type=simple
ExecStart=%h/.vcp/bin/fd-rdd
Restart=on-failure
RestartSec=5
Environment=RUST_LOG=info
Environment=XDG_RUNTIME_DIR=/run/user/%U

# Resource limits
CPUQuota=80%
MemoryMax=512M

[Install]
WantedBy=default.target