# scripts/fd-rdd.service
[Unit]
Description=fd-rdd: RDD-based elastic file indexer (user service)
After=local-fs.target
Documentation=https://github.com/AegisCompute/fd-rdd

[Service]
Type=simple
# 请确保 ~/.config/fd-rdd/config.toml 已配置 root 等参数
# 或使用 ExecStart=%h/.vcp/bin/fd-rdd --root /home/%u --snapshot-path %h/.local/share/fd-rdd/index.db
ExecStart=%h/.vcp/bin/fd-rdd
Restart=on-failure
RestartSec=5
Environment=RUST_LOG=info

# Resource limits
CPUQuota=80%
MemoryMax=512M

[Install]
WantedBy=default.target